# etc/80-usb-backup.rules
KERNEL=="sd?1", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", ENV{ID_FS_LABEL}=="USBBACKUP*", RUN+="/usr/local/bin/udev-backup.sh"
