# etc/80-usb-backup.rules
KERNEL=="sd?1", ENV{SUBSYSTEM}=="block", ENV{ID_PATH}=="*usb*", ENV{ID_FS_LABEL}=="USBBACKUP*", RUN+="/usr/local/bin/udev-backup.sh"
